# .docker/php/supervisor/queue_listen.conf
[supervisord]
nodaemon=false

[program:queue-listen]
command=php artisan queue:listen
autostart=true
autorestart=true